# ez004716-lmde3/For-Ctrip-File/rc.local
#!/bin/bash
# THIS FILE IS ADDED FOR COMPATIBILITY PURPOSES
#
# It is highly advisable to create own systemd services or udev rules
# to run scripts during boot instead of using this file.
#
# In contrast to previous versions due to parallel execution during boot
# this script will NOT be run after all other services.
#
# Please note that you must run 'chmod +x /etc/rc.d/rc.local' to ensure
# that this script will be executed during boot.

touch /var/lock/subsys/local

route add -net 10.5.115.0 netmask 255.255.255.0 gw 10.248.193.248
route add -net 10.8.153.0 netmask 255.255.255.0 gw 10.248.193.248
route add -net 10.10.10.186 netmask 255.255.255.255 gw 10.248.193.248
route add -net 10.15.220.0 netmask 255.255.255.0 gw 10.248.193.248
route add -net 10.28.97.0 netmask 255.255.255.0 gw 10.248.193.248
route add -net 10.32.56.0 netmask 255.255.255.0 gw 10.248.193.248
route add -net 10.32.57.0 netmask 255.255.255.0 gw 10.248.193.248
route add -net 10.33.128.0 netmask 255.255.255.0 gw 10.248.193.248

# 2019-0318
route add -net 10.2.59.0 netmask 255.255.255.0 gw 10.248.193.248
route add -net 10.2.75.0 netmask 255.255.255.0 gw 10.248.193.248

# 2019-0327
route add -net 10.5.121.0 netmask 255.255.255.0 gw 10.248.193.248

# 2019-0429
route add -net 10.14.84.0 netmask 255.255.255.0 gw 10.248.193.248

# 2019-0515
route add -net 10.14.45.0 netmask 255.255.255.0 gw 10.248.193.248

# 2019-0805
route add -net 10.32.116.0 netmask 255.255.255.0 gw 10.248.193.248
